# debian/stustapay.stustapay-administration-api.service
[Unit]
Description=StuStaPay Administration API
After=postgresql.service

[Service]
ExecStartPre=/usr/bin/chown -R root:stustapay /etc/stustapay; /usr/bin/chmod u=rw,g=r,o= /etc/stustapay/config.yaml
ExecStart=/usr/lib/stustapay/venv/bin/stustapay administration-api
User=stustapay
Group=stustapay

[Install]
WantedBy=multi-user.target
